# sh/sierra-launcher.service
[Unit]
Description=Sierra Launcher Daemon
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/sierra-launcher-daemon
Restart=on-failure
RestartSec=2s

# Clean up socket on exit
ExecStopPost=/bin/rm -f %t/sierra-launcher.sock

[Install]
WantedBy=graphical-session.target
